# installers/macos/com.condenser.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.condenser</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/share/condenser/condenser</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/usr/local/share/condenser</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/usr/local/share/condenser/condenser.log</string>

  <key>StandardErrorPath</key>
  <string>/usr/local/share/condenser/condenser.log</string>
</dict>
</plist>
